# openwrt/package/mosquitto/files/mosquitto.init
#!/bin/sh /etc/rc.common
# Basic init script for mosquitto
# April 2014, OpenWrt.org

START=80
APP=`which mosquitto`
CONF=/etc/mosquitto/mosquitto.conf

SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1

start() {
	fmosq
	service_start $APP -c $CONF
}

stop() {
	service_stop $APP
}
