# openwrt/package/luci/contrib/package/freifunk-policyrouting/files/etc/hotplug.d/.svn/all-wcprops
K 25
svn:wc:ra_dav:version-url
V 97
/!svn/ver/9651/luci/branches/luci-0.11/contrib/package/freifunk-policyrouting/files/etc/hotplug.d
END
